# packer/files/set-swap.service
[Unit]
Description=Swap setup service

[Service]
ExecStart=/usr/bin/scylla_swap_setup -I --swap-size SWAP_SIZE

[Install]
WantedBy=multi-user.target